# DB/vagrant_stuff/50-vagrant-mount.rules
SUBSYSTEM=="bdi",ACTION=="add",RUN+="/bin/bash /root/.udev-mount-restart-services.sh"